# rootfs/etc/s6-overlay/s6-rc.d/deluge-web-log/run
#!/command/with-contenv narch
exec s6-log -bp -- T 1
