# rootfs/etc/s6-overlay/s6-rc.d/deluged/run
#!/command/with-contenv sh
DELUGE_CONFIG_DIR=${DELUGE_CONFIG_DIR:-/config}
exec s6-setuidgid deluge deluged -d -c "$DELUGE_CONFIG_DIR" -L info 2>&1
